# configs/mywebapp.socket
[Unit]
Description=Socket Activation for MyWebApp

[Socket]

ListenStream=127.0.0.1:5200
KeepAlive=true

[Install]
WantedBy=sockets.target