# configs/mywebapp.service
[Unit]
Description=My Web App Docker Container
After=docker.service
Requires=docker.service

[Service]
TimeoutStartSec=0
Restart=always

ExecStartPre=-/usr/bin/docker stop mywebapp
ExecStartPre=-/usr/bin/docker rm mywebapp

ExecStartPre=/usr/bin/docker pull ghcr.io/taitami/mywebapp:stable

ExecStart=/usr/bin/docker run --name mywebapp --network mywebapp_net -p 5200:5200 ghcr.io/taitami/mywebapp:stable
ExecStop=/usr/bin/docker stop mywebapp

[Install]
WantedBy=multi-user.target