# configs/mywebapp.service
[Unit]
Description=MyWebApp Task Tracker Service
After=network.target mariadb.service
Requires=mywebapp.socket

[Service]
Type=simple

User=app
WorkingDirectory=/var/www/mywebapp

ExecStartPre=/usr/bin/node dist/migrate.js

ExecStart=/usr/bin/node dist/index.js --port=5200 --dbhost=127.0.0.1 --dbuser=app --dbpass=app_password --dbname=mywebapp_db
Restart=always

[Install]
WantedBy=multi-user.target